# configs/systemd/ezyshield.service
[Unit]
Description=EzyShield — adaptive IP ban daemon
Documentation=https://github.com/evertramos/ezy-shield
After=network.target ezyshield-enforcer.service
Requires=ezyshield-enforcer.service

[Service]
Type=simple
User=ezyshield
Group=ezyshield
# Load the AI API key (and any other secret) from /etc/ezyshield/.env. The
# file is written by `ezyshield init` with 0600 root:ezyshield permissions.
# Leading `-` means "don't fail the unit if the file is missing" — a fresh
# install without AI configured has no .env, and the daemon fails closed
# on AI paths only when a Resolve() call actually needs a token.
EnvironmentFile=-/etc/ezyshield/.env
ExecStart=/usr/local/bin/ezyshield run
Restart=on-failure
RestartSec=5s

# Hardening — the daemon is unprivileged; all firewall mutations go through
# the ezyshield-enforcer helper (AGENTS.md §3 / SECURITY-REVIEW.md §3).
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallFilter=@system-service
SystemCallArchitectures=native

# State and socket directories
RuntimeDirectory=ezyshield
RuntimeDirectoryMode=0755
StateDirectory=ezyshield
StateDirectoryMode=0750
LogsDirectory=ezyshield
LogsDirectoryMode=0750

[Install]
WantedBy=multi-user.target
