# configs/systemd/ezyshield-enforcer.service
[Unit]
Description=EzyShield enforcer — privileged nftables helper
Documentation=https://github.com/evertramos/ezy-shield
After=network.target

[Service]
Type=simple
User=root
# Group=ezyshield makes the process egid 'ezyshield', so the control socket (and
# its RuntimeDirectory) is created root:ezyshield with no chown needed — the
# daemon (User=ezyshield) can then connect. This avoids granting the privileged
# enforcer CAP_CHOWN just to fix socket ownership (issue #6). The group is
# created by `ezyshield init` (useradd) before this unit is installed.
Group=ezyshield
ExecStart=/usr/local/bin/ezyshield-enforcer --socket /run/ezyshield-enforcer/enforcer.sock
Restart=on-failure
RestartSec=3s

# This service holds CAP_NET_ADMIN to modify nftables; nothing else.
# No network listeners — control is via unix socket only (AGENTS.md §2).
# AmbientCapabilities grants CAP_NET_ADMIN to the process without running as
# root in the traditional sense; CapabilityBoundingSet prevents acquiring more
# (notably NOT CAP_CHOWN — socket ownership is handled via Group= above).
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=no
ProtectKernelModules=no
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX AF_NETLINK
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallFilter=@system-service
SystemCallArchitectures=native

RuntimeDirectory=ezyshield-enforcer
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
